# rc.local
#!/bin/bash

sh /usr/local/bin/start_nexres.sh
